# hw_13/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:mlflow]
command=mlflow server
    --host 0.0.0.0
    --port 5000
    --backend-store-uri sqlite:////data/mlflow.db
    --artifacts-destination /data/artifacts
    --serve-artifacts
    --workers 1
environment=MLFLOW_ENABLE_ASYNC_LOGGING="false",MLFLOW_ENABLE_ARTIFACTS_PROGRESS_TRACKING="false"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:classifier]
command=uvicorn classifier_service:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:orchestrator]
command=uvicorn orchestrator_service:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
